# director/supervisord.conf
[program:director]
command=/usr/local/bin/uwsgi /srv/stencila/hub/director/uwsgi.ini
